# packaging/libprivilege-control.manifest
<manifest>
    <request>
        <domain name="_" />
    </request>
    <assign>
        <filesystem path="/usr/bin/api_feature_loader" exec_label="none"/>
        <filesystem path="/usr/bin/slp-su" exec_label="none"/>
        <filesystem path="/etc/rc.d/init.d/smack_default_labeling" exec_label="none"/>
        <filesystem path="/etc/init.d/smack_default_labeling" exec_label="none"/>
    </assign>
</manifest>
